# configs/automatic-backup/wikijs-deployment-backup.service
[Unit]
Description=WikiJS Deployment Backup

[Service]
ExecStart=/usr/bin/wikijs-deployment-backup.sh
